# meta.timer
[Unit]
Description=Timer for Meta Clients handler

[Timer]
OnCalendar=*-*-* 0,12:00:00
Persistent=true

[Install]
WantedBy=timers.target